# dot_config/systemd/user/cliphist-text.service
[Unit]
Description=Clipboard history daemon for text (cliphist)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/wl-paste --type text --watch /usr/bin/cliphist store
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
